# cms/blog-cms.service
[Unit]
Description=Blog CMS - Obsidian to Blog Publisher
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/blog/cms
ExecStart=/root/.local/bin/uv run cms.py
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
